# ep-4-ocp-custom-kernel/openshiftcustomkernel.path
[Unit]
Description=Monitor the file openshiftcustomkernel.conf for changes

[Path]
PathChanged=/etc/openshiftcustomkerneld/openshiftcustomkernel.conf
Unit=openshiftcustomkernel.service

[Install]
WantedBy=multi-user.target